# etc/systemd/system/python_cgi.service
[Unit]
Description=Python CGI HTTP server
After=network.target

[Service]
User=webuser
WorkingDirectory=/www/html
ExecStart=/usr/bin/python3 -m http.server 8082 --bind 127.0.0.1 --cgi
Restart=on-failure

[Install]
WantedBy=multi-user.target
